# systemd/user/polkit.service
[Unit]
Description=Polkit Agent application
Documentation=man: polkit(8)

PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/libexec/polkit-gnome-authentication-agent-1 

[Install]
WantedBy=sway-session.target